# script/service/nf-nms.target
[Unit]
Description=NF NMS Platform

[Install]
WantedBy=multi-user.target
